# templates/rsnapshot-hourly.timer.j2
[Unit]
Description=rsnapshot_{{ postfix  }} hourly backup

[Timer]
# Run hourly
OnCalendar=*-*-* *:00:00
Unit=rsnapshot_{{ postfix  }}@hourly.service

[Install]
WantedBy=timers.target
